# meta-agl-demo-ces2024-ic/recipes-demo/mominjection/files/mominjection.service
[Unit]
Description=Mominjection service
After=multi-user.target

[Service]
Type=forking
ExecStart=/usr/bin/mominjection.sh

[Install]
WantedBy=multi-user.target
