# meta-agl-demo-ces2024-ic/recipes-core/cluster-service/cluster-service/cluster-service.service
[Unit]
Description=AGL Cluster Service
After=basic.target
Requires=basic.target

[Service]
Type=notify
ExecStart=/usr/bin/cluster_service

[Install]
WantedBy=basic.target
