# meta-agl-demo-ces2024-ic/recipes-core/cluster-service/cluster-service/vxcan-activate.service
[Unit]
Description=VXCAN activate Service
After=basic.target
Before=cluster-service.service
Requires=basic.target

[Service]
Type=forking
ExecStart=/usr/bin/vxcan-activate.sh

[Install]
WantedBy=basic.target
